# lxns/templates/container-default.socket
[Unit]
Description=The SSH socket of Container ${name}
[Socket]
ListenStream=${port}
