# lxns/templates/container-default.service
[Unit]
Description=Container ${name} managed by lxns
[Service]
ExecStart=/usr/bin/systemd-nspawn -jbD /var/lib/machines/${name}
KillMode=process
